# packages/jehon-service-music/usr/lib/systemd/system/mopidy.service.d/jehon-service-music.conf

[Service]
Environment="PYTHONPATH=/var/lib/python/jehon-service-music"

[Unit]
RequiresMountsFor=/mnt/cloud/musiques
