# packages/🔒_unused/jehon-service-p2p/lib/systemd/system/transmission-daemon.service.d/override.conf
[Unit]
Requires=mnt-usb1.mount
After=mnt-usb1.mount
