# packages/jehon-service-stack/lib/systemd/system/jehon-stack.service

[Unit]
Description=Update stack
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/jehon-stack
