# packages/jehon-service-kiosk/usr/lib/systemd/system/jehon-kiosk-mode.service

[Unit]
Description=Select the kiosk mode

[Service]
Type=simple
ExecStart=/usr/sbin/jh-kiosk-mode

[Install]
WantedBy=multi-user.target
