# packages/jehon-service-kiosk/usr/lib/systemd/system/jehon-kiosk-select-pictures.service

[Unit]
Description=Refresh the picture's selection
RequiresMountsFor=/mnt/cloud/photos

[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart jehon-tty@tty2.service

[Install]
WantedBy=multi-user.target
