# packages/jehon-service-kiosk/usr/lib/systemd/system/jehon-kiosk-refresh-fire.service

[Unit]
Description=Refresh the fire
RequiresMountsFor=/mnt/cloud/videos

[Service]
Type=simple
ExecStart=/usr/sbin/jehon-kiosk-fire-sync
