# packages/jehon-service-kiosk/usr/lib/systemd/system/jehon-kiosk-fire.service

[Unit]
Description=Refresh the fire [do not enable - it will be started by jehon-kiosk-mode]
RequiresMountsFor=/mnt/cloud/videos

[Service]
Type=simple
ExecStart=/usr/sbin/jehon-kiosk-fire

[Install]
WantedBy=multi-user.target
