# packages/jehon-service-database/usr/lib/systemd/system/jehon-database-backup-mariadb.timer

[Unit]
Requires=mariadb.service
RequiresMountsFor=/var/backups/

[Timer]
OnBootSec=2min
# In local time
# Check with systemd-analyze calendar "xxx"
OnCalendar=2:00
Persistent=true

[Install]
WantedBy=timers.target
