# packages/jehon-service-database/usr/lib/systemd/system/jehon-database-backup-mariadb.service

[Unit]
Requires=mariadb.service
RequiresMountsFor=/var/backups/

[Service]
Type=simple
ExecStart=/usr/sbin/jh-backup-mariadb
