# packages/jehon-service-dev/usr/lib/systemd/system/jehon-dev-user-hourly.service
[Unit]
RequiresMountsFor=/var/backups/

[Service]
Type=simple
User=jehon
ExecStart=/usr/bin/jh-dev-user-run hourly
