# packages/jehon-service-dev/usr/lib/systemd/system/jehon-dev-user-daily.service
[Unit]
RequiresMountsFor=/var/backups/

[Service]
Type=simple
User=jehon
ExecStart=/usr/sbin/jehon-dev-user-run daily

[Install]
WantedBy=multi-user.target
