# packages/jehon-service-music/usr/lib/systemd/system/mpd.service.d/jehon-service-music.conf
[Unit]
RequiresMountsFor=/mnt/cloud/musiques
