# packages/jehon-service-dev/usr/share/jehon-service-dev/etc/network/if-up.d/jh-location
#!/usr/bin/env bash

set -o errexit

# shellcheck source-path=SCRIPTDIR/../../../../../../../jehon/usr/bin
. jh-lib

exec >& >(tee -a /var/log/location.log )

date

export DEBUG=1

jh_debug "PRJ: ${JH_PKG_FOLDER}"
