# packages/jehon/usr/lib/systemd/system/keyboard-setup.service.d/🔒_jehon-override.conf
[Unit]
#
# This create a loop!!!
#
# https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=846256
#

# After=systemd-tmpfiles-setup.service
