# packages/jehon-service-music/usr/lib/systemd/system/jehon-music-update-mopidy.service

[Unit]
Description=Refresh the music database
Requires=mopidy.service
After=mopidy.service

[Service]
Type=simple
User=mopidy
ExecStart=/usr/bin/mopidy --config /etc/mopidy/mopidy.conf local scan
Environment="PYTHONPATH=/var/lib/python/jehon-service-music"

[Install]
WantedBy=multi-user.target
