# packages/jehon-service-music/usr/lib/systemd/system/jehon-music-update-mopidy.timer

[Unit]
Description=Refresh the music database

[Timer]
# OnBootSec=0
OnUnitActiveSec=1w
Persistent=true

[Install]
WantedBy=timers.target
