# packages/🔒_unused/jehon-service-p2p/usr/share/jehon-service-p2p/etc/cron.daily/jh-p2p-upload
#!/usr/bin/env bash

set -o errexit

# shellcheck source=SCRIPTDIR/../../../../../../jehon/usr/bin/jh-lib
. jh-lib

jh_log_or_tty "/var/log/jh-p2p-upload.log"

jh_exclusive "$1"

REMOTE="pcloud:Serveurs/p2p/downloaded/"

uploadOneFile() {
    FILE="$1"
    DIR="$(dirname "$FILE")"
    echo "**** $(date) Uploading '$FILE' (to $REMOTE/$DIR)"

    chmod a+w "$DIR"
    chmod a+w "$FILE"

    if [ "$DIR" != "." ]; then
        /usr/bin/rclone --config /etc/jehon/rclone.conf mkdir "${REMOTE}/${DIR}"
    fi

    /usr/bin/rclone --config /etc/jehon/rclone.conf moveto "$FILE" "${REMOTE}/${FILE}"
    echo "**** $(date) Uploading '$FILE' done"
}

cd /var/lib/transmission-daemon/downloaded/
while read -r file; do
    echo "Found $file"
    uploadOneFile "$file"
done < <(find "." -type f -printf '%P\n')
