# packages/jehon-service-headless/usr/lib/systemd/system/jehon-ssh-get-latest-key.timer

# https://wiki.archlinux.org/title/systemd/Timers
# https://manpages.ubuntu.com/manpages/jammy/en/man5/systemd.timer.5.html

[Unit]
Description=Update ssh key

[Timer]
OnBootSec=0
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
