# packages/jehon-service-kiosk/usr/lib/systemd/system/getty@tty1.service.d/jehon-kiosk-override.conf
# Inspiration from https://github.com/RPi-Distro/raspi-config/blob/master/raspi-config#L1413

[Service]
#
# Clean previous value
#
#  see:
#   $ systemctl cat getty@tty1.service
#
ExecStart=
StandardError=journal
StandardOutput=journal
RestartSec=0

# No Clear (come from /etc/systemd/system/getty@tty1.service.d/noclear.conf)
TTYVTDisallocate=no

#
# Launch script
#

# -noclear: to keep screen with infos
# ExecStart=-/sbin/agetty --autologin jehon-tty --noclear %I $TERM
Type=simple
# 1: tty number
ExecStart=/usr/sbin/jehon-kiosk-pictures 1
