# packages/jehon/usr/lib/systemd/system/jehon-cloud@.service
[Unit]
Description=JH Cloud mount "%I"
After=jehon-update-rclone.service

[Service]
Type=notify
ExecStart=/usr/sbin/jehon-cloud-mount '%I'

[Install]
WantedBy=multi-user.target
