# packages/jehon/usr/lib/systemd/system/jehon-cloud@.service
[Unit]
Description=JH Cloud mount %I
After=jehon-update-rclone.service

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount "cloud:/%i" "/mnt/cloud/$(%i)" --config=/etc/jehon/rclone.conf --ro

[Install]
WantedBy=multi-user.target
