# packages/jehon/usr/lib/systemd/system/jehon-backup-snapshot.service

[Unit]
RequiresMountsFor=/var/backups/

[Service]
Type=simple
ExecStart=/usr/bin/jh-backup-take-snapshot.sh
