# packages/jehon/usr/lib/systemd/system/jehon-docker-clean.service

[Unit]
Description=Clean-up docker
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/jh-docker-clean
