# packages/jehon/usr/lib/systemd/system/jehon-ssh-get-latest.service

[Unit]
Description=Update ssh latest key
After=network-online.target
ConditionFileIsExecutable=/usr/bin/curl

[Service]
# Latest key is the root key but published otherwise...
Type=oneshot
ExecStart=/usr/bin/curl -fsSl -o /etc/ssh/authorized_keys/latest https://jehon.github.io/infrastructure/latest-key.pub
