# packages/jehon/usr/lib/systemd/system/jehon-backup-snapshot-daily.service

[Unit]
RequiresMountsFor=/var/backups/

[Service]
Type=simple
ExecStart=/usr/sbin/jehon-backup-take-snapshot.sh daily

[Install]
WantedBy=multi-user.target
