# packages/jehon/usr/lib/systemd/system/jehon-update-rclone.service

[Unit]
Description=Update rclone
After=network-online.target

[Service]
# the VT is cleared by TTYVTDisallocate
Type=oneshot
ExecStart=/usr/bin/rclone self-update
