# packages/jehon/usr/lib/systemd/system/jehon-ssh-get-latest.service

[Unit]
Description=Update ssh latest key
After=network-online.target
ConditionFileIsExecutable=/usr/bin/curl

[Service]
# the VT is cleared by TTYVTDisallocate
Type=oneshot
ExecStart=/usr/bin/curl -fsSl -o /etc/ssh/authorized_keys/latest https://jehon.github.io/packages/latest-key.pub
