# packages/jehon/usr/lib/systemd/system/jehon-backup-snapshot-monthly.service

[Unit]
RequiresMountsFor=/var/backups/

[Service]
Type=simple
ExecStart=/usr/sbin/jh-backup-take-snapshot.sh monthly

[Install]
WantedBy=multi-user.target
