# packages/jehon/usr/lib/systemd/system/jehon-ssh-get-latest.timer

# https://wiki.archlinux.org/title/systemd/Timers
# https://manpages.ubuntu.com/manpages/jammy/en/man5/systemd.timer.5.html

[Unit]
Description=Weekly update ssh latest key

[Timer]
OnBootSec=0
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
