# packages/jehon/usr/lib/systemd/system/jehon-docker-clean.service

[Unit]
Description=Update rclone
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/jh-docker-clean
