# packages/jehon/usr/lib/systemd/system/jehon-update-rclone.timer

# https://wiki.archlinux.org/title/systemd/Timers
# https://manpages.ubuntu.com/manpages/jammy/en/man5/systemd.timer.5.html

[Unit]
Description=Weekly update rclone

[Timer]
OnBootSec=0
OnUnitActiveSec=1w
Persistent=true

[Install]
WantedBy=timers.target
