# packages/jehon/usr/lib/systemd/system/jehon-backup-snapshot.timer

[Unit]
RequiresMountsFor=/var/backups/

[Timer]
OnBootSec=2min
# In local time
# Check with systemd-analyze calendar "xxx"
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
