# packages/jehon/usr/lib/systemd/system/mnt-cloud-photos.mount
[Unit]
Description=JH Cloud mount Musiques

[Mount]
Type=rclone
What=cloud:/Photos
Where=/mnt/cloud/photos
Options=_netdev,ro,vv,nofail,noexec,allow-other,config=/etc/jehon/rclone.conf,dir-perms=0775

[Install]
WantedBy=remote-fs.target
