# packages/jehon/usr/lib/systemd/system/mnt-cloud-videos.mount
[Unit]
Description=JH Cloud mount Musiques
After=jehon-update-rclone.service

[Mount]
Type=rclone
What=cloud:/Videos
Where=/mnt/cloud/videos
Options=_netdev,ro,nofail,noexec,config=/etc/jehon/rclone.conf

[Install]
WantedBy=multi-user.target
