# packages/jehon/usr/lib/systemd/system/mnt-hardware.mount
[Unit]
Description=JH Cloud mount Musiques
After=dev-disk-by\\x2dlabel-hardware.device
Before=systemd-cryptsetup@home.service
Before=systemd-cryptsetup@swap.service

[Mount]
What=/dev/disk/by-label/hardware
Where=/mnt/hardware
Type=ext4
Options=nofail

[Install]
WantedBy=multi-user.target
